# build/ca_serve.conf
[program:generate_ca_certs]
command=/root/generate_CAcerts.sh
redirect_stderr=true
startsecs = 0
autorestart=false
startretries=0
# redirect std to docker stdout/err
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
#event to start servers

[program:cfssl_serve_production]
command=/go/bin/cfssl serve -ca=/DATA/intermediate/production/ca-production.pem
   -ca-key=/DATA/intermediate/production/ca-production-key.pem -address=0.0.0.0 -port=8889
   -dbconfig=/DATA/intermediate/production/certdb-production.json
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cfssl_serve_development]
command=/go/bin/cfssl serve -ca=/DATA/intermediate/development/ca-development.pem
   -ca-key=/DATA/intermediate/development/ca-development-key.pem -address=0.0.0.0 -port=8888
   -dbconfig=/DATA/intermediate/development/certdb-development.json
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
