# build/ca_serve.conf
[program:generate_ca_certs]
command=/root/generate_CAcerts.sh
redirect_stderr=true
startsecs = 0
autorestart=false
startretries=0
# redirect std to docker stdout/err
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
#event to start servers

[program:cfssl_serve_CAI1]
command=/go/bin/cfssl serve -address=0.0.0.0 -port %(ENV_CAI1_PORT)s
 -ca /DATA/intermediate/%(ENV_CAI1_NAME)s/ca-%(ENV_CAI1_NAME)s-2nd-full.pem
 -ca-key /DATA/intermediate/%(ENV_CAI1_NAME)s/ca-%(ENV_CAI1_NAME)s-2nd-full-key.pem
 -db-config /DATA/intermediate/%(ENV_CAI1_NAME)s/certdb-%(ENV_CAI1_NAME)s.json
 -responder /DATA/intermediate/%(ENV_CAI1_NAME)s/%(ENV_CAI1_NAME)s-ocsp.pem
 -responder-key /DATA/intermediate/%(ENV_CAI1_NAME)s/%(ENV_CAI1_NAME)s-ocsp-key.pem
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cfssl_ocspresponder_CAI1]
command=cfssl ocspserve -port=%(ENV_CAI1_OCSP_PORT)s -responses=/DATA/intermediate/%(ENV_CAI1_NAME)s/ocspdump.txt  -loglevel=%(ENV_OCSP_LOGLEVEL)s
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cfssl_serve_CAI2]
command=/go/bin/cfssl serve -address=0.0.0.0 -port %(ENV_CAI2_PORT)s
 -ca /DATA/intermediate/%(ENV_CAI2_NAME)s/ca-%(ENV_CAI2_NAME)s-2nd-noserver.pem
 -ca-key /DATA/intermediate/%(ENV_CAI2_NAME)s/ca-%(ENV_CAI2_NAME)s-2nd-noserver-key.pem
 -db-config /DATA/intermediate/%(ENV_CAI2_NAME)s/certdb-%(ENV_CAI2_NAME)s.json
 -responder /DATA/intermediate/%(ENV_CAI2_NAME)s/%(ENV_CAI2_NAME)s-ocsp.pem
 -responder-key /DATA/intermediate/%(ENV_CAI2_NAME)s/%(ENV_CAI2_NAME)s-ocsp-key.pem
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cfssl_ocspresponder_CAI2]
command=cfssl ocspserve -port=%(ENV_CAI2_OCSP_PORT)s -responses=/DATA/intermediate/%(ENV_CAI2_NAME)s/ocspdump.txt  -loglevel=%(ENV_OCSP_LOGLEVEL)s
autorestart=true
autostart=false
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0